# systemd/instagram-bot.service
[Unit]
Description=Instagram Downloader Telegram Bot
After=network.target postgresql.service redis.service

[Service]
Type=simple
WorkingDirectory=/opt/instagram-downloader-bot
EnvironmentFile=/opt/instagram-downloader-bot/.env
ExecStart=/opt/instagram-downloader-bot/instagram-bot
Restart=always
RestartSec=5
User=instagrambot

[Install]
WantedBy=multi-user.target
